# ops/systemd/kuku-v2-infer.service.d/10-operator-research-preview-gate.conf
[Service]
ExecStartPre=/usr/bin/python3 /mnt/donto-data/workspace/mobtranslate.com/training/translation/check_model_bound_release_gate.py \
    --gate /etc/mobtranslate/kuku-v24.3-operator-research-preview-gate.json \
    --gate-type operator_research_preview \
    --allowed-decision OPERATOR_RESEARCH_PREVIEW_ALLOWED \
    --model-id v24.3-joint-lexeme-dose29-s3598-20260715 \
    --artifact base_model=/mnt/donto-data/donto-resources/research/translation-training/kuku-yalanji-runpod-2026-06-30/runpod/v21.2-claude-balanced-replay-20260711T050900Z/models/v21.2-claude-balanced-replay-gvn-3epoch-lr2e-5/merged/model.safetensors \
    --artifact adapter_model=/mnt/donto-data/donto-resources/research/translation-training/kuku-yalanji-runpod-2026-06-30/runpod/v24.3-joint-lexeme-dose29-rerun1-20260715/selected-adapter/adapter_model.safetensors \
    --artifact adapter_config=/mnt/donto-data/donto-resources/research/translation-training/kuku-yalanji-runpod-2026-06-30/runpod/v24.3-joint-lexeme-dose29-rerun1-20260715/selected-adapter/adapter_config.json \
    --artifact adapter_tokenizer=/mnt/donto-data/donto-resources/research/translation-training/kuku-yalanji-runpod-2026-06-30/runpod/v24.3-joint-lexeme-dose29-rerun1-20260715/selected-adapter/tokenizer.json
