# ops/systemd/kuku-v2-infer.service
[Unit]
# Kuku Yalanji v24.3 operator-requested research preview (English -> gvn).
# Backs the homepage Translate pane and /labs/v2. Loopback-only, no secrets.
Description=Kuku Yalanji v24.3 unverified research-preview inference
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ajax
Group=ajax
WorkingDirectory=/mnt/donto-data/workspace/mobtranslate.com
Environment=PYTHONUNBUFFERED=1
Environment=OMP_NUM_THREADS=4
Environment=MKL_NUM_THREADS=4
Environment=HF_HOME=/tmp/huggingface
Environment=HF_HUB_OFFLINE=1
Environment=TRANSFORMERS_OFFLINE=1
Environment=TOKENIZERS_PARALLELISM=false
# Citer-leak scar: cap glibc arenas so RSS stays bounded under MemorySwapMax=0.
Environment=MALLOC_ARENA_MAX=2
ExecStart=/mnt/donto-data/donto-resources/research/translation-training/kuku-yalanji-runpod-2026-06-30/venvs/v14-cpu-py312/bin/python \
    /mnt/donto-data/workspace/mobtranslate.com/training/translation/serve_v2_infer.py \
    --model-dir /mnt/donto-data/donto-resources/research/translation-training/kuku-yalanji-runpod-2026-06-30/runpod/v21.2-claude-balanced-replay-20260711T050900Z/models/v21.2-claude-balanced-replay-gvn-3epoch-lr2e-5/merged \
    --adapter-dir /mnt/donto-data/donto-resources/research/translation-training/kuku-yalanji-runpod-2026-06-30/runpod/v24.3-joint-lexeme-dose29-rerun1-20260715/selected-adapter \
    --adapter-sha256 dd61583a60df2d538989e963e104cb626d78965d300e4e473e9a82ef59c04502 \
    --model-id v24.3-joint-lexeme-dose29-s3598-20260715 \
    --host 127.0.0.1 --port 7955 \
    --dtype float32 --threads 4 \
    --max-chars 400 --max-new-tokens 208 --num-beams 1 \
    --no-repeat-ngram-size 4 --repetition-penalty 1.10 --length-penalty 1.0 \
    --max-waiting 2 --max-time 45
# /health only answers once the model has loaded AND warmed up.
# Match TimeoutStartSec: allow 450 two-second probes when the mounted drive is busy.
ExecStartPost=/bin/bash -c 'for attempt in {1..450}; do /usr/bin/curl -fsS http://127.0.0.1:7955/health >/dev/null 2>&1 && exit 0; sleep 2; done; exit 1'
Restart=always
RestartSec=10
# Loading 2.7 GB of weights is normally quick, but shared-drive research scans can
# saturate /mnt/donto-data. Preserve forward progress instead of restart-looping.
TimeoutStartSec=900
TimeoutStopSec=30
# User-facing latency path -> normal scheduling priority (NOT deprioritized).
# The merged model is ~1.37B params (larger than the assumed 600M base): fp32 is
# The dynamic v24.3 adapter measured 7.1G settled RSS and 8.003G peak cgroup
# memory. Leave explicit startup headroom while preventing host-wide growth.
MemoryHigh=8G
MemoryMax=9G
MemorySwapMax=0
# Loose CPU backstop so a burst can't monopolize the shared box.
CPUQuota=400%
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
