# ops/systemd/kuku-v2-infer.service
[Unit]
# Kuku Yalanji v21.2 research-preview translator (English -> gvn).
# Backs the mobtranslate.com /labs/v2 page. Loopback-only, no secrets.
Description=Kuku Yalanji v21.2 research-preview translation inference (/labs/v2)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ajax
Group=ajax
WorkingDirectory=/mnt/donto-data/workspace/mobtranslate.com
Environment=PYTHONUNBUFFERED=1
Environment=OMP_NUM_THREADS=4
Environment=MKL_NUM_THREADS=4
Environment=HF_HOME=/tmp/huggingface
Environment=HF_HUB_OFFLINE=1
Environment=TRANSFORMERS_OFFLINE=1
Environment=TOKENIZERS_PARALLELISM=false
# Citer-leak scar: cap glibc arenas so RSS stays bounded under MemorySwapMax=0.
Environment=MALLOC_ARENA_MAX=2
ExecStart=/mnt/donto-data/donto-resources/research/translation-training/kuku-yalanji-runpod-2026-06-30/venvs/v14-cpu-py312/bin/python \
    /mnt/donto-data/workspace/mobtranslate.com/training/translation/serve_v2_infer.py \
    --model-dir /mnt/donto-data/donto-resources/research/translation-training/kuku-yalanji-runpod-2026-06-30/runpod/v21.2-claude-balanced-replay-20260711T050900Z/models/v21.2-claude-balanced-replay-gvn-3epoch-lr2e-5/merged \
    --model-id v21.2-claude-balanced-replay \
    --host 127.0.0.1 --port 7955 \
    --dtype float32 --threads 4 \
    --max-chars 400 --max-new-tokens 208 --num-beams 1 \
    --max-waiting 8 --max-time 45
# /health only answers once the model has loaded AND warmed up.
ExecStartPost=/bin/bash -c 'for attempt in {1..180}; do /usr/bin/curl -fsS http://127.0.0.1:7955/health >/dev/null && exit 0; sleep 2; done; exit 1'
Restart=always
RestartSec=10
TimeoutStartSec=420
TimeoutStopSec=30
# User-facing latency path -> normal scheduling priority (NOT deprioritized).
# The merged model is ~1.37B params (larger than the assumed 600M base): fp32 is
# ~5.3G resident with an ~8G transient 2x load peak, so the cap matches the proven
# sibling (mobtranslate-translate-v2 = 8G). MemorySwapMax=0 + MALLOC_ARENA_MAX=2
# keep memory bounded and non-swapping per the citer-leak scar.
MemoryHigh=7G
MemoryMax=8G
MemorySwapMax=0
# Loose CPU backstop so a burst can't monopolize the shared box.
CPUQuota=600%
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
