# ops/systemd/mobtranslate-hybrid-warm.timer
[Unit]
Description=Verify and warm enabled MobTranslate adapters every 15 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=false

[Install]
WantedBy=timers.target
