# ops/systemd/mobtranslate-web.service
[Unit]
Description=mobtranslate.com immutable Next.js release
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ajax
Group=ajax
WorkingDirectory=/mnt/donto-data/mobtranslate/current/runtime/apps/web
EnvironmentFile=/opt/mobtranslate/web.env
EnvironmentFile=-/mnt/donto-data/mobtranslate/current/release-runtime.env
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3300
Environment=MOBTRANSLATE_RUNTIME_PUBLIC_DIR=/mnt/donto-data/workspace/mobtranslate.com/apps/web/public
Environment=MOBTRANSLATE_REPO_ROOT=/mnt/donto-data/workspace/mobtranslate.com
Environment=MOBTRANSLATE_DICTIONARIES_ROOT=/mnt/donto-data/workspace/mobtranslate.com/dictionaries
ExecStart=/usr/bin/node /mnt/donto-data/mobtranslate/current/runtime/apps/web/server.js
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
KillSignal=SIGTERM
SuccessExitStatus=143
LimitNOFILE=65536
UMask=0027
NoNewPrivileges=true
CapabilityBoundingSet=
LockPersonality=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
ReadWritePaths=/mnt/donto-data/mobtranslate-storage
ReadWritePaths=/mnt/donto-data/workspace/mobtranslate.com/apps/web/public/word-img
ReadWritePaths=/mnt/donto-data/workspace/mobtranslate.com/apps/web/public/wotd
RemoveIPC=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
