# ops/systemd/mobtranslate-translate-v2.service.d/10-dictionary-cpu-gate.conf
[Service]
ExecStartPre=/usr/bin/python3 /mnt/donto-data/workspace/mobtranslate.com/training/translation/check_dictionary_cpu_deployment_gate.py --gate /etc/mobtranslate/kuku-cpu-deployment-gate.json --model-id v21.1-codex-synthetic-direct-gvn-3epoch-lr2e-5 --model-file /mnt/donto-data/donto-resources/research/translation-training/kuku-yalanji-runpod-2026-06-30/runpod/v21.1-codex-synthetic-direct-20260710T171317Z/models/v21.1-codex-synthetic-direct-gvn-3epoch-lr2e-5/merged/model.safetensors
