# ops/systemd/mobtranslate-translate-v2.service.d/20-sentence-generation-gate.conf
[Service]
ExecStartPre=/usr/bin/python3 /mnt/donto-data/workspace/mobtranslate.com/training/translation/check_model_bound_release_gate.py --gate /etc/mobtranslate/kuku-v21.1-sentence-generation-gate.json --gate-type model_sentence_generation --allowed-decision SENTENCE_GENERATION_ALLOWED --model-id v21.1-codex-synthetic-direct-gvn-3epoch-lr2e-5 --model-file /mnt/donto-data/donto-resources/research/translation-training/kuku-yalanji-runpod-2026-06-30/runpod/v21.1-codex-synthetic-direct-20260710T171317Z/models/v21.1-codex-synthetic-direct-gvn-3epoch-lr2e-5/merged/model.safetensors
